# mon0-setup.service
# mon0-setup.service — systemd unit to create a Wi-Fi monitor interface.
#
# Started two ways:
#   1. At boot via "systemctl enable" (script waits up to WAIT_SECS for the adapter).
#   2. On hot-plug via the udev rule (99-mon0-setup.rules) — adapter already present.
# Type=oneshot + RemainAfterExit=yes prevents double-runs if both triggers fire.
#
# Install instructions (run as root):
#   cp setup_mon0.sh       /usr/local/sbin/
#   chmod +x               /usr/local/sbin/setup_mon0.sh
#   cp mon0-setup.service  /etc/systemd/system/
#   cp 99-mon0-setup.rules /etc/udev/rules.d/
#   systemctl daemon-reload
#   udevadm control --reload-rules
#   systemctl enable --now mon0-setup.service
#
# Logs:  journalctl -u mon0-setup.service

[Unit]
Description=Create Wi-Fi monitor interface (mon0) for packet capture
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Timeout must exceed WAIT_SECS in setup_mon0.sh (default 40s) plus setup time
TimeoutStartSec=60
ExecStart=/usr/local/sbin/setup_mon0.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
